# scripts/open_apps_fullscreen/com.user.open-apps-fullscreen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.open-apps-fullscreen</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>sleep 30 &amp;&amp; /Users/vaibhavisingh/bin/open_apps_fullscreen</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/open-apps-fullscreen.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/open-apps-fullscreen.error.log</string>
</dict>
</plist>